# person/novelScratch/deploy/ubuntu/mq/kafka/kafka.service
[Unit]
Requires=zookeeper.service
After=zookeeper.service

[Service]
Type=simple
User=chenye
ExecStart=/bin/sh -c '/home/chenye/kafka/bin/kafka-server-start.sh /home/chenye/kafka/config/server.properties > /home/chenye/kafka/kafka.log 2>&1'
ExecStop=/home/chenye/kafka/bin/kafka-server-stop.sh
Restart=on-abnormal

[Install]
WantedBy=multi-user.target